# deploy/clav-backup.timer
[Unit]
Description=Run clav-backup nightly

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
